# roles/sway/files/on_hdmi_connected.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/bin/sway-laptop-display.sh"
SUBSYSTEM=="power_supply", KERNEL=="LID", ACTION=="change", RUN+="/usr/local/bin/sway-laptop-display.sh"
