# roles/gotify-listener/templates/gotify-listen.service.j2
[Unit]
Description=Gotify desktop notification listener
After=graphical-session.target
PartOf=graphical-session.target

[Service]
ExecStart={{ gotify_listener_bin_dir }}/gotify-listen
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
