# roles/sway-extra/files/on_hdmi_connected.rules
SUBSYSTEM=="drm", RUN+="/bin/bash /usr/local/etc/sway/configure-laptop-display.sh"
